# supervisor.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[include]
files=/etc/supervisor/conf.d/*.conf

[program:env]
command=/usr/bin/env
startsecs=0
autorestart=false

[program:master]
# Directory where your sitr's project files are located
directory=/srv/

command=daphne -b 0.0.0.0 -p 8000 -t 1800 --porxy-headers config.asgi:application

# Automatically start and recover processes
autostart=true
autorestart=true

# Choose where you want your log to go
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:uwsgi]
# Directory where your sitr's project files are located
directory=/srv/

command=uwsgi server.ini -t 1800 --protocol=http

# signal used to kill process (default TERM)
stopsignal=QUIT

# Automatically start and recover processes
autostart=true
autorestart=true

# Choose where you want your log to go
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

